# deploy/fips-restart.service
[Unit]
Description=FIPS Daemon Restart (for recovery testing)
After=fips.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart fips
